# config/systemd/user/redshift.service.d/local-config.conf
[Service]
ExecStart=
ExecStart=/usr/bin/redshift -c /home/daniel/.config/redshift.conf
