# config/systemd/user/wallpaper.service
[Unit]
Description=Symlink a background/wallpaer image to $HOME/.cache/wallpaper.jpg

[Service]
Type=oneshot
ExecStart=%h/bin/wallpaper
